# deploy/com.btcbot.live.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.btcbot.live</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOURUSER/.local/bin/btcbot</string>
    <string>live</string>
    <string>--config</string>
    <string>/Users/YOURUSER/for_claude/config/default.yaml</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOURUSER/for_claude</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/YOURUSER/for_claude/logs/btcbot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOURUSER/for_claude/logs/btcbot_err.log</string>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
